# raspberry-pi/rpi_scanner.service
[Unit]
Description=RTL-SDR RF Anomaly Scanner with Bluetooth
After=network.target bluetooth.target

[Service]
Type=simple
User=test
SupplementaryGroups=dialout bluetooth
WorkingDirectory=/home/test/trackapp/raspberry-pi
EnvironmentFile=-/etc/default/rpi_scanner
Environment=BT_TRANSPORT=ble
Environment=BLE_NAME_PREFIX=TrackScan
ExecStart=/home/test/.venvs/rfscan/bin/python /home/test/trackapp/raspberry-pi/rpi_live_scanner.py scanner --bt-transport ${BT_TRANSPORT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
