# raspberry-pi/rpi_scanner.service
[Unit]
Description=RTL-SDR RF Anomaly Scanner with Bluetooth
After=network.target bluetooth.target

[Service]
Type=simple
User=test
SupplementaryGroups=dialout bluetooth
PermissionsStartOnly=true
WorkingDirectory=/home/test/trackapp/raspberry-pi
EnvironmentFile=-/etc/default/rpi_scanner
ExecStartPre=/bin/sh -c '/usr/bin/rfcomm release /dev/rfcomm0 || true'
ExecStartPre=/bin/sh -c 'if [ -n "${BT_TARGET_MAC}" ]; then /usr/bin/rfcomm bind /dev/rfcomm0 "${BT_TARGET_MAC}"; fi'
ExecStart=/home/test/.venvs/rfscan/bin/python /home/test/trackapp/raspberry-pi/rpi_live_scanner.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
